# supervisord.conf
[supervisord]
user = root

[program:rsyslog]
command = /sbin/rsyslogd -n -c5 -f /etc/rsyslog.conf
priority = 10

[program:cat]
process_name = cat
command = cat /dev/consolelog
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
priority = 20

[program:cron]
command = /usr/sbin/crond -n
priority = 30

[program:sqlgrey]
process_name = sqlgrey
command = /usr/sbin/sqlgrey
priority = 100

[program:postfix]
process_name = master
directory = /etc/postfix
command = /postfix.sh
directory = /etc/postfix

[program:amavisd]
command = /usr/sbin/amavisd -c /etc/amavisd/amavisd.conf foreground

[program:clamd]
command	= /clamd.amavisd.sh

[program:dovecot]
command	= /usr/sbin/dovecot -c /etc/dovecot/dovecot.conf -F
